# debian/ksm.postinst
#!/bin/sh
echo "Enable  at boot time"
sudo update-rc.d ksm defaults
